# workers/default.ini
[program:default-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/lifeapp-api/artisan queue:work redis --sleep=30 --tries=1 --timeout=600 --memory=4096 --max-time=900
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=apache
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/html/lifeapp-api/workers/logs/worker.log
stopwaitsecs=900
